# pkg/DEBIAN/prerm
#!/bin/bash

systemctl stop werearud
systemctl disable werearud
