# pkg/etc/systemd/system/werearud.service
[Unit]
Description=Where are you server

[Service]
ExecStart=/usr/bin/werearud

[Install]
WantedBy=multi-user.target
